# buttons.service
[Unit]
Description=Arcade Button Control

[Service]
ExecStart=/home/pi/busyboard/button-control.py

[Install]
WantedBy=default.target

